# packages/bfgminer/bfgminer.service
[Unit]
Description=Bitcoin miner featuring overclocking, monitoring, fan speed control and remote management. For FPGA/GPU/CPU Bitcoin mining.

[Service]
ExecStart=/usr/bin/bfgminer -T --config /etc/bfgminer.conf --syslog

[Install]
WantedBy=network-online.target
